# deploys/systemd/system/klever-celery-development.service
[Unit]
Description=Klever Celery Development Service
After=network.target local-fs.target remote-fs.target

[Service]
Type=simple
EnvironmentFile=/etc/default/klever
ExecStartPre=/bin/sh -c "cd ${KLEVER_DEPLOYMENT_DIRECTORY}/klever/bridge"
ExecStart=/usr/local/bin/celery beat -A bridge worker --loglevel=info

[Install]
WantedBy=multi-user.target
